# supervisor/elasticsearch.conf
[program:elasticsearch]
command=/opt/elasticsearch/bin/elasticsearch -f
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/elasticsearch.out.log
stderr_logfile=/var/log/supervisor/elasticsearch.err.log

