# supervisor/kibana.conf
[program:kibana]
command=/opt/kibana/bin/kibana
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/kibana.out.log
stderr_logfile=/var/log/supervisor/kibana.err.log

