# supervisor/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D -d
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/sshd.our.log
stderr_logfile=/var/log/supervisor/sshd.err.log
